# etc/supervisord.d/logrotate.conf
[program:logrotate]
command=/opt/logrotate.sh
autorestart=true
stdout_logfile=/var/log/supervisor/logrotate.log
redirect_stderr=true
